# Linux/Systemd/examples/ddns-update/update_ddns.service
# /etc/systemd/system/update_ip.service
[Unit]
Description=Update DDNS and log IP

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update_ddns.sh
